# deploy/runtime/systemd/llm-nginx.service
[Unit]
Description=nginx reverse proxy for llm.nalits.com
After=network-online.target llmapi.service
Wants=network-online.target
Requires=llmapi.service

[Service]
Type=simple
ExecStartPre=/usr/sbin/nginx -t -c /home/ubuntu/llmapi/deploy/runtime/nginx/nginx.conf
ExecStart=/usr/sbin/nginx -c /home/ubuntu/llmapi/deploy/runtime/nginx/nginx.conf -g "daemon off;"
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
Slice=llm.nalits.slice

[Install]
WantedBy=multi-user.target
