# deploy/runtime/systemd/llm.nalits.slice
[Unit]
Description=Memory limit for llm.nalits.com (llmapi + nginx)

[Slice]
MemoryAccounting=yes
MemoryMax=400M
MemorySwapMax=0
